# root/etc/systemd/system/dhcpcd@eth0.service.d/override.conf
[Unit]
OnFailure=hostapd.service

[Service]
Restart=always
