# root/etc/systemd/system/wpa_supplicant@wlan0.service.d/override.conf
[Unit]
After=hostapd.service access-point.service