# root/usr/lib/systemd/system/rslsync.service.d/override.conf
[Service]
Group=thebox
User=thebox
ExecStart=
ExecStart=/usr/bin/rslsync --nodaemon --config /home/thebox/.config/rslsync/rslsync.conf
