# root/etc/systemd/system/minidlna.service.d/override.conf
[Unit]
Conflicts=minidlna-rebuild.service minidlna-rescan.service