# root/etc/systemd/system/hostapd.service.d/override.conf
[Unit]
Before=wpa_supplicant@wlan0.service
After=systemd-resolved.service systemd-networkd.service access-point.service dhcpcd@eth0.service