# root/usr/lib/systemd/system/hostapd.service.d/override.conf
[Unit]
After=network.target dnsmasq_ap.service supplicant@wlan0.service
BindsTo=dnsmasq_ap.service

[Service]
ExecStartPre=/usr/bin/sh -c 'CHANNEL=`iw dev wlan0 info|grep channel|cut -d" " -f2`; sed -i "s/^channel=.*/channel=${CHANNEL:-6}/" /etc/hostapd/hostapd.conf'
