# root/etc/systemd/system/minidlna-rescan.service
[Unit]
Description=minidlna server
After=network.target
Conflicts=minidlna.service minidlna-rebuild.service

[Service]
Type=simple
User=minidlna
Group=minidlna
ExecStart=/usr/bin/minidlnad -S -r
ProtectSystem=full
ProtectHome=on
PrivateDevices=on
NoNewPrivileges=on

[Install]
WantedBy=multi-user.target